# stub-runner.jar
#!/bin/bash
#
#    .   ____          _            __ _ _
#   /\\ / ___'_ __ _ _(_)_ __  __ _ \ \ \ \
#  ( ( )\___ | '_ | '_| | '_ \/ _` | \ \ \ \
#   \\/  ___)| |_)| | | | | || (_| |  ) ) ) )
#    '  |____| .__|_| |_|_| |_\__, | / / / /
#   =========|_|==============|___/=/_/_/_/
#   :: Spring Boot Startup Script ::
#

### BEGIN INIT INFO
# Provides:          spring-cloud-contract-stub-runner-boot
# Required-Start:    $remote_fs $syslog $network
# Required-Stop:     $remote_fs $syslog $network
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Spring Cloud Contract Stub Runner Boot
# Description:       Spring Cloud Contract Stub Runner Boot
# chkconfig:         2345 99 01
### END INIT INFO

[[ -n "$DEBUG" ]] && set -x

# Initialize variables that cannot be provided by a .conf file
WORKING_DIR="$(pwd)"
# shellcheck disable=SC2153
[[ -n "$JARFILE" ]] && jarfile="$JARFILE"
[[ -n "$APP_NAME" ]] && identity="$APP_NAME"

# Follow symlinks to find the real jar and detect init.d script
cd "$(dirname "$0")" || exit 1
[[ -z "$jarfile" ]] && jarfile=$(pwd)/$(basename "$0")
while [[ -L "$jarfile" ]]; do
  if [[ "$jarfile" =~ init\.d ]]; then
    init_script=$(basename "$jarfile")
  else
    configfile="${jarfile%.*}.conf"
    # shellcheck source=/dev/null
    [[ -r ${configfile} ]] && source "${configfile}"
  fi
  jarfile=$(readlink "$jarfile")
  cd "$(dirname "$jarfile")" || exit 1
  jarfile=$(pwd)/$(basename "$jarfile")
done
jarfolder="$( (cd "$(dirname "$jarfile")" && pwd -P) )"
cd "$WORKING_DIR" || exit 1

# Inline script specified in build properties


# Source any config file
configfile="$(basename "${jarfile%.*}.conf")"

# Initialize CONF_FOLDER location defaulting to jarfolder
[[ -z "$CONF_FOLDER" ]] && CONF_FOLDER="${jarfolder}"
# shellcheck source=/dev/null
[[ -r "${CONF_FOLDER}/${configfile}" ]] && source "${CONF_FOLDER}/${configfile}"

# ANSI Colors
echoRed() { echo $'\e[0;31m'"$1"$'\e[0m'; }
echoGreen() { echo $'\e[0;32m'"$1"$'\e[0m'; }
echoYellow() { echo $'\e[0;33m'"$1"$'\e[0m'; }

# Initialize PID/LOG locations if they weren't provided by the config file
[[ -z "$PID_FOLDER" ]] && PID_FOLDER="/var/run"
[[ -z "$LOG_FOLDER" ]] && LOG_FOLDER="/var/log"
! [[ "$PID_FOLDER" == /* ]] && PID_FOLDER="$(dirname "$jarfile")"/"$PID_FOLDER"
! [[ "$LOG_FOLDER" == /* ]] && LOG_FOLDER="$(dirname "$jarfile")"/"$LOG_FOLDER"
! [[ -x "$PID_FOLDER" ]] && echoYellow "PID_FOLDER $PID_FOLDER does not exist. Falling back to /tmp" && PID_FOLDER="/tmp"
! [[ -x "$LOG_FOLDER" ]] && echoYellow "LOG_FOLDER $LOG_FOLDER does not exist. Falling back to /tmp" && LOG_FOLDER="/tmp"

# Set up defaults
[[ -z "$MODE" ]] && MODE="auto" # modes are "auto", "service" or "run"
[[ -z "$USE_START_STOP_DAEMON" ]] && USE_START_STOP_DAEMON="true"

# Create an identity for log/pid files
if [[ -z "$identity" ]]; then
  if [[ -n "$init_script" ]]; then
    identity="${init_script}"
  else
    identity=$(basename "${jarfile%.*}")_${jarfolder//\//}
  fi
fi

# Initialize log file name if not provided by the config file
[[ -z "$LOG_FILENAME" ]] && LOG_FILENAME="${identity}.log"

# Initialize stop wait time if not provided by the config file
[[ -z "$STOP_WAIT_TIME" ]] && STOP_WAIT_TIME="60"

# Utility functions
checkPermissions() {
  touch "$pid_file" &> /dev/null || { echoRed "Operation not permitted (cannot access pid file)"; return 4; }
  touch "$log_file" &> /dev/null || { echoRed "Operation not permitted (cannot access log file)"; return 4; }
}

isRunning() {
  ps -p "$1" &> /dev/null
}

await_file() {
  end=$(date +%s)
  let "end+=10"
  while [[ ! -s "$1" ]]
  do
    now=$(date +%s)
    if [[ $now -ge $end ]]; then
      break
    fi
    sleep 1
  done
}

# Determine the script mode
action="run"
if [[ "$MODE" == "auto" && -n "$init_script" ]] || [[ "$MODE" == "service" ]]; then
  action="$1"
  shift
fi

# Build the pid and log filenames
PID_FOLDER="$PID_FOLDER/${identity}"
pid_file="$PID_FOLDER/${identity}.pid"
log_file="$LOG_FOLDER/$LOG_FILENAME"

# Determine the user to run as if we are root
# shellcheck disable=SC2012
[[ $(id -u) == "0" ]] && run_user=$(ls -ld "$jarfile" | awk '{print $3}')

# Ensure the user actually exists
id -u "$run_user" &> /dev/null || unset run_user

# Run as user specified in RUN_AS_USER
if [[ -n "$RUN_AS_USER" ]]; then
    if ! [[ "$action" =~ ^(status|run)$ ]]; then
        id -u "$RUN_AS_USER" || {
            echoRed "Cannot run as '$RUN_AS_USER': no such user"
            exit 2
        }
        [[ $(id -u) == 0 ]] || {
            echoRed "Cannot run as '$RUN_AS_USER': current user is not root"
            exit 4
        }
    fi
    run_user="$RUN_AS_USER"
fi

# Issue a warning if the application will run as root
[[ $(id -u ${run_user}) == "0" ]] && { echoYellow "Application is running as root (UID 0). This is considered insecure."; }

# Find Java
if [[ -n "$JAVA_HOME" ]] && [[ -x "$JAVA_HOME/bin/java" ]]; then
    javaexe="$JAVA_HOME/bin/java"
elif type -p java > /dev/null 2>&1; then
    javaexe=$(type -p java)
elif [[ -x "/usr/bin/java" ]];  then
    javaexe="/usr/bin/java"
else
    echo "Unable to find Java"
    exit 1
fi

arguments=(-Dsun.misc.URLClassPath.disableJarChecking=true $JAVA_OPTS -jar "$jarfile" $RUN_ARGS "$@")

# Action functions
start() {
  if [[ -f "$pid_file" ]]; then
    pid=$(cat "$pid_file")
    isRunning "$pid" && { echoYellow "Already running [$pid]"; return 0; }
  fi
  do_start "$@"
}

do_start() {
  working_dir=$(dirname "$jarfile")
  pushd "$working_dir" > /dev/null
  if [[ ! -e "$PID_FOLDER" ]]; then
    mkdir -p "$PID_FOLDER" &> /dev/null
    if [[ -n "$run_user" ]]; then
      chown "$run_user" "$PID_FOLDER"
    fi
  fi
  if [[ ! -e "$log_file" ]]; then
    touch "$log_file" &> /dev/null
    if [[ -n "$run_user" ]]; then
      chown "$run_user" "$log_file"
    fi
  fi
  if [[ -n "$run_user" ]]; then
    checkPermissions || return $?
    if [ $USE_START_STOP_DAEMON = true ] && type start-stop-daemon > /dev/null 2>&1; then
      start-stop-daemon --start --quiet \
        --chuid "$run_user" \
        --name "$identity" \
        --make-pidfile --pidfile "$pid_file" \
        --background --no-close \
        --startas "$javaexe" \
        --chdir "$working_dir" \
        -- "${arguments[@]}" \
        >> "$log_file" 2>&1
      await_file "$pid_file"
    else
      su -s /bin/sh -c "$javaexe $(printf "\"%s\" " "${arguments[@]}") >> \"$log_file\" 2>&1 & echo \$!" "$run_user" > "$pid_file"
    fi
    pid=$(cat "$pid_file")
  else
    checkPermissions || return $?
    "$javaexe" "${arguments[@]}" >> "$log_file" 2>&1 &
    pid=$!
    disown $pid
    echo "$pid" > "$pid_file"
  fi
  [[ -z $pid ]] && { echoRed "Failed to start"; return 1; }
  echoGreen "Started [$pid]"
}

stop() {
  working_dir=$(dirname "$jarfile")
  pushd "$working_dir" > /dev/null
  [[ -f $pid_file ]] || { echoYellow "Not running (pidfile not found)"; return 0; }
  pid=$(cat "$pid_file")
  isRunning "$pid" || { echoYellow "Not running (process ${pid}). Removing stale pid file."; rm -f "$pid_file"; return 0; }
  do_stop "$pid" "$pid_file"
}

do_stop() {
  kill "$1" &> /dev/null || { echoRed "Unable to kill process $1"; return 1; }
  for ((i = 1; i <= STOP_WAIT_TIME; i++)); do
    isRunning "$1" || { echoGreen "Stopped [$1]"; rm -f "$2"; return 0; }
    [[ $i -eq STOP_WAIT_TIME/2 ]] && kill "$1" &> /dev/null
    sleep 1
  done
  echoRed "Unable to kill process $1";
  return 1;
}

force_stop() {
  [[ -f $pid_file ]] || { echoYellow "Not running (pidfile not found)"; return 0; }
  pid=$(cat "$pid_file")
  isRunning "$pid" || { echoYellow "Not running (process ${pid}). Removing stale pid file."; rm -f "$pid_file"; return 0; }
  do_force_stop "$pid" "$pid_file"
}

do_force_stop() {
  kill -9 "$1" &> /dev/null || { echoRed "Unable to kill process $1"; return 1; }
  for ((i = 1; i <= STOP_WAIT_TIME; i++)); do
    isRunning "$1" || { echoGreen "Stopped [$1]"; rm -f "$2"; return 0; }
    [[ $i -eq STOP_WAIT_TIME/2 ]] && kill -9 "$1" &> /dev/null
    sleep 1
  done
  echoRed "Unable to kill process $1";
  return 1;
}

restart() {
  stop && start
}

force_reload() {
  working_dir=$(dirname "$jarfile")
  pushd "$working_dir" > /dev/null
  [[ -f $pid_file ]] || { echoRed "Not running (pidfile not found)"; return 7; }
  pid=$(cat "$pid_file")
  rm -f "$pid_file"
  isRunning "$pid" || { echoRed "Not running (process ${pid} not found)"; return 7; }
  do_stop "$pid" "$pid_file"
  do_start
}

status() {
  working_dir=$(dirname "$jarfile")
  pushd "$working_dir" > /dev/null
  [[ -f "$pid_file" ]] || { echoRed "Not running"; return 3; }
  pid=$(cat "$pid_file")
  isRunning "$pid" || { echoRed "Not running (process ${pid} not found)"; return 1; }
  echoGreen "Running [$pid]"
  return 0
}

run() {
  pushd "$(dirname "$jarfile")" > /dev/null
  "$javaexe" "${arguments[@]}"
  result=$?
  popd > /dev/null
  return "$result"
}

# Call the appropriate action function
case "$action" in
start)
  start "$@"; exit $?;;
stop)
  stop "$@"; exit $?;;
force-stop)
  force_stop "$@"; exit $?;;
restart)
  restart "$@"; exit $?;;
force-reload)
  force_reload "$@"; exit $?;;
status)
  status "$@"; exit $?;;
run)
  run "$@"; exit $?;;
*)
  echo "Usage: $0 {start|stop|force-stop|restart|force-reload|status|run}"; exit 1;
esac

exit 0
PK  V            	  META-INF/   PK           PK  V               META-INF/MANIFEST.MF}j0þDJ}}PОkȒYI	}v(-=-(gNxA
ƻ䆳PEԢ|+QWtP?l`Vn9+Z"	x8YE8hn"O*">B__r#YU!iaٟHxtǜ!WIHMSqUrEE?~fyTm劵fe%s0gPKlpn     PK
    J@V               org/PK
    J@V               org/springframework/PK
    J@V               org/springframework/boot/PK
    J@V                org/springframework/boot/loader/PK
    J@V            (   org/springframework/boot/loader/wrapper/PK
    L@V               META-INF/maven/PK
    L@V            5   META-INF/maven/org.springframework.boot.experimental/PK
    L@V            N   META-INF/maven/org.springframework.boot.experimental/spring-boot-thin-wrapper/PK  J@V            <   org/springframework/boot/loader/wrapper/ThinJarWrapper.classZx\U׽tӔn4aM6lki5IK-elݰA
"*`iJ(("`іϹwo6ɶ{Ϝsf8[Dtcc<N^^.3<'{)Sq,WeU:uR+x3Tg
q<Β:yy,$f9s5>KS<yc&P//Ey,^*_ѸAލ^ne:/Bun+djΫDj\s\,5X+sqEBcm:ѹCuܥsT˼$n,w}b:'tb|4N˻Wtޤf/ï-^WQ/U\uڥ5_+M6NI1^|#Z&Q:,|?tU/'ZtM(E1_篈鿪9Cyw^:#_P^\7t8:?wt+؀[_JSK7/[pYxS&wִxlD<ӫ±64,\ܚ5/^j!S,[\fp21խl_ɫ_ySqcݪF/Xne
{DVoO2IMGc5m:nld⋙\hk?r*'qx{C$"!%±UdTNvBrIfLz,o$F)JۜdgMGu$ݑ͉ƚD&G5Қ]pr9-g֞fJ+>mI66{2`i	f\-ғL9x4=KwE%x*ZNĢq4> E0U&*51>TWsB͚Pa,gtOO8C#*v8\M,\HO뒑ToLTit<Z,!NnEv-d{4N>T~RV!ga>NuE7A6DSiLNB3A	̳Fǣ;&edM5NmW 6e7dMFr#ooH+gˇ<4HG,Җ1YD;a$:}vA)8@>D=M`tLo! f@(M'[O  Û"qu%q,{S|*V`scD*ͣ3&+/[Ҵ|W&cPj2;=
#(V&TlU 7lU1X3DDjbVڝ+#%li``X+4xF(yPGCYq=vȖ̥#IiI.vDo@ˠ %pL
ڴM5ާh0ߒMES28SWd/?`t\ci$k~!?e}=y|gP"rR`4Ƞoѷse| 4~?l/L|ԠiS?g8 %X3CyJGl hLp<HrHYPh/!__!Lac
moF$̿3vH]2j{_?d=a
`Ijj~hЛr~Ort8$~x0bi }=
H6ed%,V'AJ/-Mu[Kw7+IA?6W̤Td3rW̨V70n[CeLzՑ-I%9+2UO`PdK-4zzֵ,\]2*e?uJLoKpX,|[L g3~)JVF4|G IU9Oe283uRÙZ)a[C'l6jl/iz.vxm0AUlBH.f߅^PcCjgEo
G0!WzQB@a6"~yep,*Nc*Q1\ցz%b$[6Hs~QҬ+(((eDY*
g$T"(WeE1q<֓edgFl	*i@JsUS[Z"t#DSes574gt]c.օXkDn0Q7+l^;r2zSlrXV!U?=όU!Ұ.HIB*֛].S@N+u`vvy̎PE(Y>;XټdL,L{j6tlrT:s7,ping=ySS´4_}\"NY P ջ!fe˒|މ\Ze'Q[~PRj)݊Ϝz-qd2]@[nR
@:Tґutl>F<N_%PգApۂl6谮$	nv6}RCܑG]$?v^|$#	"@ȸR+bNf<̈Hcs?rFui(W]00k;;Ϝy\nn%Rx :?XbԲƵ*ʇ\OK"6'D䗠"WaeȤJHIQa?>'~>6rz$	Hx}S;*=@hܮ94]~AL._3_R{3j`cDؠ-ݣT:EZ{ֺEyĩn~`ح:f䧯݁-g|w挝490zx-Gb94G0ޙ3n_\7nxEc F:ĭ>r4Vr2s?[+3ɐ!=C}P^2VWamу9t<#tR)uxR53i#C1ڍe<cIgzTm2ЫOٕJ!	ѣ*L70;Eiʽ4*F0CTUedir>BT*ɍg\C%t9U* #%#g0PY G?M%'hF/H|"HSi@EF7uֺJ]F]rݾ34f;3wCn|US24>{3349ງvwF&TOl:Mpf-pVGWz"t:(8&qfV8Xr2Qbr Of~ZLZ[́~l34g}rOoE&r@ܕ5
NMAgɭi
9Mקe*QiC.X YrՄ1]}4%CV? ],k 멉n@
1\SySWZrO(L axk|$h.pmIؾǩ"$vhsN n|g!,*]t[Vl^sH}?1kl҅rlӢMBVYOr*>z4$|8Csrw  ѼF/YgOȻN};hb.&KXpֺ|aUEGEƵR8 &i~~m'T+%wKeN삓ӏwg#!-+@hIa8 pgi27#|@iNRDs6ZH;h)}N$л PRߠ@y	w&ѷvh'wx.=p'=̗#|5Ƿӣ|chփv3rB8Yr2z42nR!i4
U}10`AE)`SD#':N,c1;&tԪ&a{|awaBl;ԏ!0HO6ر>:>͇VAܷeMUngS꾛&]
\UA.tmd_jAf+NCT[w-@Yja.ft[	>Z%L5asի}tbkv9p.$g =GpYg\I?u"/Q_Bf|`}^VBUc42t]H?S9-yʋ$:NH9͐&g*UyѪQe(1iWф/8!?͛ArT]VÔdu.iBL7M].R_8Cևt;	a3Ԇ zvc<;i!:ӤJiHSBڍӅ`:jG5хx!Vὖ.Eu0]Px]qt5w]6~B'l*p:g
_ƽ	^oC;nMZRК@W[6xzMEG$ށAwwC/&RQS__EMFuVi}݅ulo5zw|K^@JA _w{<?@m=*Ƌ+}:[5hG|Plj5)fy!2(ehc
sPM{0<[|44]Tn_[閭G#48PR
6of̞j4YbjPB.bΦ}d9S'>NoNİ%$K[GHA\	qi[Q"C8ʞ1NE1mȎ׳?G۰c&"!aq`dy,>,-I\{,G[K>oن[k|[]S޼nSsl{0	3]}mj#!Eoݸ}Pr]2tbiW$59GqhgGթ|t}rn#m#3i N+GKTi#zkT~bwSn{z-\\i9Gېl%*hga[>(:Vt[)Z@(ZrۨvD;LG#e`YoIP=)7qȵyG>Cĉn+٧B]=Fib1~,dY]U>#N2Z$ejX]l7#[̰T;iJ wwsP76|t7Z߉Mfރ0s=+g'W>aI#" 24ϓxME4/ټZKiVnvnN^Ay%")V^KW%tkz^F7r/~=X?q!އYykbϰ%ys`ͫ|)afy>-@G-9L<Ֆafndi*=nX+.G445.v?g
Q)ARC֨3fdkw2/K՘?ե1?^
{_Aӭv~£K۹qͷv<`16,{C*N3	nKLJ&wұY&XM]!Db~
Iוpq75qqڸGci Gp x6P0!>G,V-|݂,M^mC%;ڌz,M>J5yfi\~Tm"dvǧˏI\"F?PKp/-  /  PK  V               lib/ PK           PK  V            
   lib/.empty PK           PK
    V               org/springframework/cloud/PK
    V            #   org/springframework/cloud/contract/PK
    V            .   org/springframework/cloud/contract/stubrunner/PK
    V            5   org/springframework/cloud/contract/stubrunner/server/PK
    V            )   META-INF/maven/org.springframework.cloud/PK
    V            P   META-INF/maven/org.springframework.cloud/spring-cloud-contract-stub-runner-boot/PK  V            I   org/springframework/cloud/contract/stubrunner/server/StubRunnerBoot.classRMo1}nlZ("ni%)Ģpݸ[yũU ~bDJ>o>	16RR肿.dگVasp8aH-x5JgV7c4F;+r+WgZZ^I{I.%c@bĸma6(|=*NMB1C6&ї5$unt8nK2wTQ[r"YGL%Xie~E9VlsU14sۀ<L(zޠ][]!L<L?/IKqjj˷kM·`'ĩJhm\	7>+,IY|׵ˉ]M<sDy ߗ8yF~kx#wɶwv	O~p|5?@v@zt~#"[=!PK  v  PK  V               META-INF/thin-rabbit.properties  L0 (A
)_5y]wNMkA!E;•3v+
ܪdsUȞ$;VPK_X      PK  V               META-INF/thin-kafka.properties	  Nt %XDEnσ;ONM{A!E;•a*Ѣ@J6˷
f1jPK]X      PK  V               application.yml-1 @ѝSXDRbl	!Rj#czhY9qďU .#z' BrATܶX,<Q; KܝfA5O(\7Y3+ʶ=^9/KuL8LPK`      PK  V            W   META-INF/maven/org.springframework.cloud/spring-cloud-contract-stub-runner-boot/pom.xmlUn0=;_aJBqP݊vUd:Q+K$')h)ɜil#-U#X'd'8KAq]	5dwy9>Erg+'&{sFr,
mgx|B~\_94,y8d$ۢ-@0Dқę$W3L;{	rULX`6=^Wc)3[sYHY8cqBe,},mEfcy \+o1GjJzntJb%$Nc7ˢd'jVomO~.;l&w
.omzѦkio6h*p
~M0*\P.Ldxc=?%\q0/MkUDcKt|0;>}xuΌx~HBJ#$P0BFF. <;SNdtPֲz$M%޶Y/u/#-2`EϺߧKu6wػINڹ$CTϭέ6[PGV3!//weͤCv{뷳yJֻ-GPK
b    PK  V            ^   META-INF/maven/org.springframework.cloud/spring-cloud-contract-stub-runner-boot/pom.properties%A =w)qx@ DmͧKt;+Cl$Z$sEub MxgPK[T   b   PK  V           	         A9$  META-INF/  PK  Vlpn                v$  META-INF/MANIFEST.MFPK
    J@V                      A%  org/PK
    J@V                      A%  org/springframework/PK
    J@V                      A&  org/springframework/boot/PK
    J@V                       A8&  org/springframework/boot/loader/PK
    J@V            (          Av&  org/springframework/boot/loader/wrapper/PK
    L@V                      A&  META-INF/maven/PK
    L@V            5          A&  META-INF/maven/org.springframework.boot.experimental/PK
    L@V            N          A<'  META-INF/maven/org.springframework.boot.experimental/spring-boot-thin-wrapper/PK  J@Vp/-  /  <           '  org/springframework/boot/loader/wrapper/ThinJarWrapper.classPK  V                     A	@  lib/PK  V           
           =@  lib/.emptyPK
    V                      Aw@  org/springframework/cloud/PK
    V            #          A@  org/springframework/cloud/contract/PK
    V            .          A@  org/springframework/cloud/contract/stubrunner/PK
    V            5          A<A  org/springframework/cloud/contract/stubrunner/server/PK
    V            )          AA  META-INF/maven/org.springframework.cloud/PK
    V            P          AA  META-INF/maven/org.springframework.cloud/spring-cloud-contract-stub-runner-boot/PK  V  v  I           DB  org/springframework/cloud/contract/stubrunner/server/StubRunnerBoot.classPK  V_X                 D  META-INF/thin-rabbit.propertiesPK  V]X                 $E  META-INF/thin-kafka.propertiesPK  V`                 E  application.ymlPK  V
b    W           F  META-INF/maven/org.springframework.cloud/spring-cloud-contract-stub-runner-boot/pom.xmlPK  V[T   b   ^           xI  META-INF/maven/org.springframework.cloud/spring-cloud-contract-stub-runner-boot/pom.propertiesPK      W  XJ    